# debian/cron.d
#
# Regular cron jobs for the dorian package
#
0 4	* * *	dorian	/opt/dorian/bin/dorian/cleanup
